# build/packaging/scripts/preremove.sh
#!/bin/bash
systemctl stop podman-volume-stratis 2>/dev/null || true
systemctl disable podman-volume-stratis 2>/dev/null || true
